# bin/ptrack.service
[Unit]
Description=Ptrack
After=network.target

[Service]
Type=idle
Restart=on-failure
User=ec2-user
ExecStart=/bin/bash -c 'cd /home/ec2-user/ptrack/ && flask --app ptrack --debug run --host=0.0.0.0'

[Install]
WantedBy=multi-user.target
